# units/mbsync-oneshot.timer.in
[Unit]
Description=mbsync query timer
ConditionPathExists=%h/.mbsyncrc

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15m

[Install]
WantedBy=default.target
